# config/org.danlynn.clamav.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!--
LaunchAgent script that executes the clamav.rb script once a day and generates
an HTML report.  This script was generated by running: clamav.rb -i [TIME]
This script will not be activated until the next restart of the computer.
Similarly, if this script is deleted (or clamav.rb -u) the running
of clamav.rb on the specified interval will continue until the next restart
of the computer.
 -->
<dict>
	<key>Label</key>
	<string>org.danlynn.clamav</string>
	<key>ProgramArguments</key>
	<array>
		<string>sh</string>
		<string><%= root_dir %>/clamav.rb</string>
	</array>
	<key>QueueDirectories</key>
	<array/>
	<key>ServiceDescription</key>
	<string>Run the clamav.rb script once a day</string>
	<key>StandardErrorPath</key>
	<string><%= root_dir %>/log/launchd_err.log</string>
	<key>StandardOutPath</key>
	<string><%= root_dir %>/log/launchd_out.log</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer><%= time.hour %></integer>
		<key>Minute</key>
		<integer><%= time.min %></integer>
	</dict>
	<key>WatchPaths</key>
	<array/>
</dict>
</plist>
