# config/org.danlynn.clamav.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
LaunchAgent script that executes the clamav.rb script once a day and generates
an HTML report.  This script was generated by running: clamav.rb -i [TIME]
Any changes to this script (add/delete/update) will not be activated until
the next time that you logout and then log back in to this account.
 -->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>org.danlynn.clamav</string>
	<key>ProgramArguments</key>
	<array>
		<string>ruby</string>
		<string><%= root_dir %>/clamav.rb</string>
		<string>--config</string>
		<string><%= root_dir %>/<%= config_path %></string>
	</array>
	<key>QueueDirectories</key>
	<array/>
	<key>ServiceDescription</key>
	<string>Run the clamav.rb script once a day</string>
	<key>StandardErrorPath</key>
	<string><%= root_dir %>/log/launchd_err.log</string>
	<key>StandardOutPath</key>
	<string><%= root_dir %>/log/launchd_out.log</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer><%= time.hour %></integer>
		<key>Minute</key>
		<integer><%= time.min %></integer>
	</dict>
	<key>WatchPaths</key>
	<array/>
</dict>
</plist>
